# vagrant/upstart-private-ledger.conf
start on vagrant-mounted
stop on runlevel [06]
stop on paysharesd-reprovision

respawn

exec /paysharesd-src/build/paysharesd --conf /paysharesd-src/vagrant/paysharesd-private-ledger.cfg --load
